# debian/postinst
#!/bin/sh
# postinst script for gerrit
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        local confirmation
        local gerrit_path=/var/lib/gerrit
        local dirs="$gerrit_path /etc/gerrit /var/log/gerrit /var/cache/gerrit"
        adduser --system --home $gerrit_path --no-create-home --group --shell /bin/bash --disabled-password gerrit
        chgrp -R gerrit $dirs
        chmod -R g+w $dirs
        su -c "ssh-keygen -q -f "$gerrit_path/etc/ssh_host_rsa_key" -N '' -t rsa" gerrit
        su -c "ssh-keygen -q -f "$gerrit_path/etc/ssh_host_dsa_key" -N '' -t dsa" gerrit
        echo
        echo "When asked to install required Bouncy Castle Crypto v144 say no, as it is already included as dependency"
        echo "Press any key to continue..."
        read confirmation
        su -c "java -jar $gerrit_path/bin/gerrit.war init -d $gerrit_path" gerrit
        cp /var/lib/gerrit/bin/gerrit.sh /etc/init.d/gerrit
        update-rc.d gerrit defaults
        service gerrit start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
