# debian/postinst
#!/bin/sh
# postinst script for gerrit
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
	local GERRIT_SITEPATH=/var/lib/gerrit
	local DIRS="$GERRIT_SITEPATH /etc/gerrit /var/log/gerrit /var/cache/gerrit"
        adduser --system --home $GERRIT_SITEPATH --no-create-home --group --shell /bin/bash --disabled-password gerrit
        chgrp -R gerrit $DIRS
        chmod -R g+w $DIRS
	su -c "ssh-keygen -q -f "$GERRIT_SITEPATH/etc/ssh_host_rsa_key" -N '' -t rsa" gerrit
	su -c "ssh-keygen -q -f "$GERRIT_SITEPATH/etc/ssh_host_dsa_key" -N '' -t dsa" gerrit
	echo "When asked to install required Bouncy Castle Crypto v144 say no, as it is already included as dependency"
	echo "Press any key to continue..."
	read || true
        su -c "java -jar $GERRIT_SITEPATH/bin/gerrit.war init -d $GERRIT_SITEPATH" gerrit
        update-rc.d gerrit defaults
        service gerrit start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
